# docker/s6-overlay/etc/services.d/clash/run
#!/usr/bin/with-contenv bash

echo "[services.d] clash: service is running..."

exec s6-setuidgid root clash -d $WORKSPACE -ext-ui $WORKSPACE/ui